# templates/etc/profile.d/hadoop.sh.j2
export PATH=$PATH:{{ hadoop_usr_dir }}/sbin:$PATH:{{ hadoop_usr_dir }}/bin